# target/linux/ramips/base-files/etc/init.d/set-irq-affinity
#!/bin/sh /etc/rc.common

START=99

start() {
	if grep -q 'processor.*: 2' /proc/cpuinfo; then
		mask=4
	elif grep -q 'processor.*: 1' /proc/cpuinfo; then
		mask=2
	else
		return
	fi

	# Spread irq load over 3 cpu cores
	echo e > /sys/class/net/eth0/queues/rx-0/rps_cpus

	for irq in $(grep "mt76..e" /proc/interrupts | cut -d: -f1 | sed 's, *,,')
	do
			echo "$mask" > "/proc/irq/$irq/smp_affinity"
			[ $mask = 4 ] && mask=8
	done

}
